# countdown.service
[Unit]
Description=Countdown server

[Service]
Type=simple
ExecStart=/home/pi/countdown/countdown.py

[Install]
WantedBy=multi-user.target
